# topvpn.service
[Unit]
Description=TopVPN auto login
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Host and user are read from ~/.topvpn.conf (sourced by the script)
# Optional debug: uncomment to enable
# Environment=TOPVPN_DEBUG=1
WorkingDirectory=/home/al/t
ExecStart=/bin/bash /home/al/t/topvpn.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
