# topvpn.service
[Unit]
Description=TopVPN auto login
Wants=network-online.target topvpnhelper.service
After=network-online.target topvpnhelper.service

[Service]
Type=simple
# Host and user are read from ~/.topvpn.conf (sourced by the script)
# Optional debug: uncomment to enable
# Environment=TOPVPN_DEBUG=1
#WorkingDirectory=/home/al/.topvpn
ExecStart=/home/al/.local/bin/topvpn.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
